# deploy/systemd/mister-morph.service
[Unit]
Description=mistermorph (agent daemon, sandboxed)
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=10
StartLimitBurst=10

[Service]
Type=simple

# Example deployment layout:
# - Binary:  /opt/morph/mistermorph
# - Config:  /opt/morph/config.yaml
# - Skills:  /opt/morph/skills (set skills.dirs accordingly in config)
# Or ExecStart=/opt/morph/mistermorph --config /opt/morph/config.yaml --log-level info telegram
ExecStart=/opt/morph/mistermorph --config /opt/morph/config.yaml --log-level info serve --server-bind 127.0.0.1 --server-port 8787

# Secrets/config via environment (recommended).
# Suggested split:
# - /opt/morph/morph.env: non-secret config (permissions: 0640)
# - /opt/morph/morph.secrets.env: secret values / secret_ref env vars (permissions: 0600)
#
# Common env vars:
# - MISTER_MORPH_LLM_API_KEY
# - MISTER_MORPH_SERVER_AUTH_TOKEN
# - Any auth_profile secret_ref env vars (e.g. JSONBILL_API_KEY)
EnvironmentFile=-/opt/morph/morph.env
EnvironmentFile=-/opt/morph/morph.secrets.env

Restart=on-failure
RestartSec=2s

# Run as a dedicated user, so you can manage ownership/ACLs explicitly.
User=morph
Group=morph

# systemd-managed writable dirs (kept minimal).
# - StateDirectory: persistent data (memory files + guard approvals state)
# - CacheDirectory: ephemeral cache (file_cache_dir / telegram downloads / write_file)
StateDirectory=morph
StateDirectoryMode=0700
CacheDirectory=morph
CacheDirectoryMode=0700

# Keep all persistent writes inside StateDirectory, and all temp/file-tool writes inside CacheDirectory.
WorkingDirectory=/var/lib/morph
Environment=MISTER_MORPH_FILE_CACHE_DIR=/var/cache/morph

# If you want to allow the agent to read a specific project directory, prefer a read-only bind mount:
# BindReadOnlyPaths=/home/you/projects:/workspace

NoNewPrivileges=yes
UMask=0077

PrivateTmp=yes
PrivateDevices=yes
DevicePolicy=closed

ProtectSystem=strict
ProtectHome=true
ProtectControlGroups=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid

RestrictSUIDSGID=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
SystemCallArchitectures=native

# Needs outbound network (OpenAI/Telegram/web_search) + local unix sockets.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

[Install]
WantedBy=multi-user.target
